# mpv_mods/Etc_Systemd_user/deck.service
[Service]
ExecStart=/bin/bash /home/pi/scripts/launchdeck.sh
